# templates/supervisor.conf
[program:mikbotam]
directory=/var/www/html/Saldo
command=php7.2 Core.php
