# templates/supervisor.conf
[program:mikbotam]
command=php7.2 /var/www/html/Saldo/Core.php
